# usb/81-thumbdrive.rules
KERNEL=="sd?1", SYMLINK+="thumbdrive", RUN+="/root/exercises/jaars/run-on-insert.sh"
